# packages/kcgroups-dmemcg/kcgroups-dmemcg-suite.install
post_install() {
    echo "==> Orchestrating VRAM: Enabling dmemcg-booster services..."
    # Enable the system-level hierarchy manager
    systemctl enable dmemcg-booster.service

    # Enable the user-level manager for all current and future users
    systemctl --global enable dmemcg-booster.service

    echo "==> VRAM protection is now active. A logout/login or reboot is recommended."
}

post_upgrade() {
    post_install
}
